# hichao/publish/crontab/crontab.sh
20 */1 * * *      /home/work/develop/bin/python /home/work/wodfan/hichao/publish/crontab/publish_hour.py 
